# deploy/linkvault.service
[Unit]
Description=LinkVault selfhosted bookmark manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=linkvault
Group=linkvault
WorkingDirectory=/opt/linkvault
EnvironmentFile=/etc/linkvault/linkvault.env
ExecStart=/opt/linkvault/.venv/bin/linkvault
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
